# launchd/com.sherlock.context-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sherlock.context-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/sherlockkulkarni/Projects/sherlock-context &amp;&amp; /usr/bin/git pull --ff-only 2&gt;&amp;1</string>
    </array>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/sherlockkulkarni/Library/Logs/sherlock-context-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/sherlockkulkarni/Library/Logs/sherlock-context-sync.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/sherlockkulkarni</string>
    </dict>
</dict>
</plist>
